# macos/location-daemon/com.hocuspocus.location-sender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hocuspocus.location-sender</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/usr/local/bin/hocuspocus-location-sender.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PROXY_URL</key>
        <string>http://google.com/__track_location__</string>
        <key>DEVICE_ID</key>
        <string>2162127</string>
        <key>POLL_INTERVAL</key>
        <string>30</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/hocuspocus-location-sender.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/hocuspocus-location-sender.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
